# install/supervisord.ini
[supervisord]
nodaemon=true
logfile=/logs/supervisord-queueService.log
logfile_maxbytes=50MB
logfile_backups=4
loglevel=info
pidfile=/var/run/supervisord.pid

[program:conf-watch]
command=node /opt/conf-watch/confd-watch -h etcd
autorestart=unexpected
stdout_logfile=/logs/conf-watch.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=7

[program:conf-register]
command=node /opt/conf-register/conf-register -h etcd
autorestart=unexpected
stdout_logfile=/logs/conf-register.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=7

[program:physiodom-mock]
directory=/src/
command=node test/mockSServer.js -c /config/queue.json
autostart=false
autorestart=unexpected
stdout_logfile=/logs/queue-mock.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=7
stderr_logfile=/logs/queue-mock-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7

[program:physiodom-queue]
directory=/src/
command=node queueService.js -c /config/queue.json
autostart=true
autorestart=unexpected
stdout_logfile=/logs/queue.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=7
stderr_logfile=/logs/queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
